# tika-gotenberg/rootfs/etc/s6-overlay/s6-rc.d/gotenberg/run
#!/command/with-contenv bashio
cd /home/gotenberg
exec s6-setuidgid gotenberg /usr/bin/gotenberg