# tika-gotenberg/rootfs/etc/s6-overlay/s6-rc.d/tika/run
#!/command/with-contenv bashio
/usr/bin/tika